# cron/crontab
# run this every 10 minutes
*/6 * * * * mk_pf_rss.sh
